# cmd/admin/supervisorconf/upspin_directory.conf
[program:upspin-dirserver]
environment=GOOGLE_APPLICATION_CREDENTIALS=/var/www/upspin-720f98729598.json,HOME=/var/www
command=/var/www/dirserver -endpoint=gcp,local -config="storage=Postgres,user=postgres,host=localhost,port=5432,password=postgres,sslmode=disable"
directory=/var/www
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/upspin/directory.err.log
stdout_logfile=/var/log/upspin/directory.out.log
user=www-data
